# ceph/init/bc-ceph-init.service
[unit]
Description=init-ceph-chown

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/bc-ceph-init

[Install]
WantedBy=multi-user.target

